# curator4-1696382071.conf
[program:curator4-1696382071]
command=python3 /home/curator4/curator4-1696382071/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-1696382071/test.err.log
stdout_logfile=/home/curator4/curator4-1696382071/test.out.log